# install/supervisor/hanabi-live.conf
[program:hanabi-live]
directory=/root/hanabi-live
command=/root/hanabi-live/hanabi-live
environment=GOPATH=/root/go
redirect_stderr=true
stdout_logfile=/root/hanabi-live/logs/hanabi-live.log
; 1 GB x 10 = 10 GB total of logs
stdout_logfile_maxbytes=1GB
stdout_logfile_backups=10
